# dockerfiles/piraeus-oci/tmpl/.service
[Unit]
Description=Piraeus ${COMPONENT} OCI Container

[Service]
Type=simple
EnvironmentFile=-/etc/default/%p
WorkingDirectory=/opt/piraeus/oci/${COMPONENT}
ExecStartPre=/bin/sh -c "/opt/piraeus/oci/runc delete -f piraeus-${COMPONENT} || true"
ExecStart=/opt/piraeus/oci/runc run --pid-file=/run/piraeus-${COMPONENT}.pid piraeus-${COMPONENT}
ExecStopPost=/opt/piraeus/oci/runc delete -f piraeus-${COMPONENT}
PIDFile=/run/piraeus-${COMPONENT}.pid
Restart=always
RestartSec=5
RestartPreventExitStatus=10
KillMode=mixed

[Install]
WantedBy=multi-user.target